# scripts/yenten-stratum.service
[Unit]
Description=Yenten Mining Pool Stratum Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/yenten-pool
ExecStart=/usr/bin/php /var/www/yenten-pool/scripts/stratum_server.php
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/yenten-pool/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
